# supervisord.conf
; supervisor config file.
; Start => /usr/bin/python3 /usr/local/bin/supervisord -c /etc/supervisord.conf

[supervisord]
logfile=/home/kudin/StravaKudos/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=true               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)
user = root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; use a unix:// URL  for a unix socket

[program:StravaKudos]
directory               = /home/kudin/StravaKudos
command                 = /home/kudin/StravaKudos/strava-kudos-amd64
process_name            = %(program_name)s_%(process_num)02d
numprocs                = 1
autostart               = true
autorestart             = true
stdout_logfile          = /home/kudin/StravaKudos/strava-kudos_info.log
stdout_logfile_maxbytes = 1MB
stderr_logfile          = /home/kudin/StravaKudos/strava-kudos_error.log
stderr_logfile_maxbytes = 1MB
user                    = root